# server/supervisor.conf
[supervisord] 
logfile = /home/docker/logs/supervisor.log
logfile_maxbytes = 10MB
childlogdir = /home/docker/logs/

[program:esproxy]
command = /home/docker/code/es_proxy/run_proxy.py

