# server/supervisor.conf
[supervisord] 
logfile = /opt/logs/supervisor.log
logfile_maxbytes = 10MB
childlogdir = /opt/logs/

[program:esproxy]
command = /opt/code/es_proxy/run_proxy.py

